# trunk/Library/LaunchAgents.disabled/zsh.fink.completion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>zsh.fink.completion</string>
	<key>OnDemand</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/init/zsh/local-functions/darwin/_FinkMakeListCache</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ServiceDescription</key>
	<string>Update the fink completion cache whenever it needs to be changed</string>
	<key>WatchPaths</key>
	<array>
		<string>/sw/var/lib/dpkg/status</string>
	</array>
</dict>
</plist>
